# roles/router/templates/dhcpcd@wan.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/dhcpcd -q -w %I -C resolv.conf
